# judson.service

[Unit]
Description= Judson

[Service]
User= root
ExecStart=/bin/judson
Restart = always

[Install]
WantedBy=multi-user.target

